# macos_launch_agents/Library/LaunchAgents/com.user.wezterm.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.user.wezterm</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/sh</string>
    <string>-c</string>
    <string>launchctl setenv WEZTERM_CONFIG_FILE $HOME/.config/wezterm/wezterm.lua</string>
  </array>
  <key>RunAtLoad</key>
  <true/>
</dict>
</plist>
